# firewall/seer-firewall.service
[Unit]
Description=SEER Firewall API Service
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/seer
ExecStart=/opt/seer/venv/bin/python /opt/seer/api.py
Restart=always
RestartSec=3
Environment="PYTHONUNBUFFERED=1"

[Install]
WantedBy=multi-user.target
